# lib/main/file/etc/update-motd.d/99-pst-banner
#!/bin/sh

figlet -c -w 80 -m0 'Welcome to PST!!!'
